# roles/nukuhiva/templates/srv-moorea.automount.j2
[Unit]
  Description=Moorea {{ item }} automount
[Automount]
  Where=/srv/moorea/{{ item }}
[Install]
  WantedBy=multi-user.target